# templates/default/init-script.erb
#!/bin/sh

# chkconfig: 2345 80 20

PERL5LIB="/var/www/<%= node['movabletype-opensource']['domain'] %>/cgi-bin/mt/extlib:/var/www/<%= node['movabletype-opensource']['domain'] %>/cgi-bin/mt/lib"

user="<%= node['movabletype-opensource']['user'] %>"
group="<%= node['movabletype-opensource']['user'] %>"
start_server="<%= node['movabletype-opensource']['perl_install_path'] %>/bin/start_server"
starman="<%= node['movabletype-opensource']['perl_install_path'] %>/bin/starman"
workers=<%= node['movabletype-opensource']['workers'] %>
port=<%= node['movabletype-opensource']['psgi_port'] %>
app="/var/www/<%= node['movabletype-opensource']['domain'] %>/cgi-bin/mt/mt.psgi"
pid_file="/var/run/movabletype"
status_file="/var/run/movabletype-status"
log="/var/www/<%= node['movabletype-opensource']['domain'] %>/logs/start_server.log"

start() {
    PERL5LIB=${PERL5LIB} \
    nohup ${start_server} \
        --port=${port} \
        --signal-on-hup=QUIT \
        --pid-file=${pid_file} \
        --status-file=${status_file} -- \
            ${starman} \
            --user=${user} \
            --group=${group} \
            --workers=${workers} \
            --max-reqs-per-child=10000 \
            --interval=1 \
            --preload-app \
            -E production \
            -a ${app} 1>${log} 2>&1 &
}

stop() {
    kill -TERM $(cat ${pid_file})
}

graceful() {
    PERL5LIB=${PERL5LIB} \
    ${start_server} \
        --restart \
        --pid-file=${pid_file} \
        --status-file=${status_file}
}

case "$1" in 
    start)
        start
        ;;
    stop)
        stop
        ;;
    restart)
        stop
        start
        ;;
    graceful)
        graceful
        ;;
    *)
        exit 2;
        ;;
esac
